# glazed-whale.service
[Unit]
Description=Glazed Whale - Franchiser Mining Bot
After=network.target

[Service]
Type=simple
User=donut-agent
WorkingDirectory=/home/donut-agent/.openclaw/workspace/glazed-whale
ExecStart=/usr/bin/node scripts/monitor.js
Restart=always
RestartSec=10
StandardOutput=append:/var/log/glazed-whale/output.log
StandardError=append:/var/log/glazed-whale/error.log

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
